# deploy/systemd/covex-backup.service
[Unit]
Description=Covex nightly SQLite backup (crash-consistent snapshot to /var/backups/covex)
Documentation=file:/mnt/HC_Volume_105579109/Covex27/deploy/backup-covex.sh
RequiresMountsFor=/mnt/HC_Volume_105579109
OnFailure=covex-alert@%n.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-covex.sh
TimeoutStartSec=1h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
