# deploy/systemd/covex-restore-drill.service
[Unit]
Description=Covex weekly backup restore drill (verifies the newest snapshot restores)
Documentation=file:/mnt/HC_Volume_105579109/Covex27/deploy/restore-drill-covex.sh
Wants=covex-backup.service
After=covex-backup.service
OnFailure=covex-alert@%n.service

[Service]
Type=oneshot
Environment=BACKUP_ROOT=/mnt/HC_Volume_105579109/covex-backups
ExecStart=/usr/local/bin/restore-drill-covex.sh
TimeoutStartSec=30min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
