# deploy/systemd/covex-alert@.service
[Unit]
Description=Covex failure alert for %i

[Service]
Type=oneshot
ExecStart=/usr/local/bin/covex-alert.sh %i
