# deploy/covex-backend.service
[Unit]
Description=Covex27 Backend - Rust Kaspa Covenant Explorer API
After=network.target

[Service]
Type=simple
User=kasparov
WorkingDirectory=/home/kasparov/Covex27
Environment="KASPA_NETWORK=testnet-12"
Environment="KASPA_WRPC_URL=ws://178.105.76.81:17110"
Environment="BIND_ADDR=0.0.0.0:3001"
Environment="DB_PATH=/home/kasparov/Covex27/covex.db"
Environment="COVENANT_TREASURY_ADDRESS=kaspatest:qpyfz03k6quxwf2jglwkhczvt758d8xrq99gl37p6h3vsqur27ltjhn68354m"
Environment="RUST_LOG=covex27_backend=debug,kaspa_wrpc=info"
ExecStart=/home/kasparov/Covex27/backend/target/release/covex27-backend
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
